# time_to_go_notifier.service
[Unit]
Description=Time To Go Notification

[Service]
ExecStart=/home/alex/projects/Time-To-Go/time_to_go.sh

[Install]
WantedBy=multi-user.target
